# dining_room/sensor_data_logger.service
[Unit]
Description=DHT22 Poller Dining Room
After=influxdb.service

[Service]
Type=simple
User=pi
WorkingDirectory= /home/pi/pi_temps/dining_room/
ExecStart=/usr/bin/python3 /home/pi/pi_temps/dining_room/sensor_data_logger.py
Restart=always

[Install]
WantedBy=multi-user.target